# s6/uranus/run
#!/command/execlineb -P
/uranus/bin/firewalld-gateway -v 5  --sql-driver=sqlite --config=/uranus/firewalld-gateway.toml